# wordpress/templates/varnish.service.erb
[Unit]
Description=Varnish HTTP accelerator
Documentation=https://www.varnish-cache.org/docs/4.1/ man:varnishd

<%
  total_ram_mb    = node['memory']['total'].gsub('kB', '').to_i / 1024
  varnish_mb      = (total_ram_mb * 0.20).to_i
  thread_pools    = 4
  thread_pool_max = [varnish_mb / thread_pools / 2, 25].max
-%>
[Service]
Type=simple
LimitNOFILE=131072
LimitMEMLOCK=82000
ExecStart=/usr/sbin/varnishd -j unix,user=vcache -F -a 0.0.0.0:80 -T 127.0.0.1:6082 -f /etc/varnish/default.vcl -S /etc/varnish/secret -s malloc,<%= varnish_mb %>m -p thread_pools=<%= thread_pools %> -p thread_pool_min=5 -p thread_pool_max=<%= thread_pool_max %> -p thread_pool_timeout=300
ExecReload=/usr/share/varnish/reload-vcl
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true

[Install]
WantedBy=multi-user.target